# osc_bridge.service
[Unit]
Description=Behringer Wing to Xilica Bridge Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /home/pi/wing_to_xilica_simple.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
