# _backup_STABLE_2026_01_20/ptz_cam.service
[Unit]
Description=PTZ Camera Control Bridge
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/behringer-mixer
ExecStart=/home/pi/behringer-mixer/venv/bin/python -u /home/pi/behringer-mixer/ptz_cam.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
